# root/etc/supervisord.d/mysql.ini
[program:mysqld]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/sbin/mysqld
autostart=true
process_name=%(program_name)s
startsecs=5
stopsignal=INT
user=mysql
redirect_stderr=true
stdout_logfile=/var/log/mysql.log
stderr_logfile=/var/log/mysql.err
