# root/etc/supervisord.d/tomcat.ini
[program:tomcat]
command=/opt/tomcat/bin/catalina.sh run
environment=JAVA_HOME="/usr/lib/jvm/java-1.8.0-openjdk-amd64/jre"
autostart = true
autorestart=true
